# supervisord.conf
[supervisord]
childlogdir=/var/log/supervisor
nodaemon=true

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server

[program:monit]
command=/usr/bin/monit
